# deploy/kejar-prestasi.service
[Unit]
Description=Kejar Prestasi (Node SSR)
Documentation=https://github.com/your-org/kejarprestasi
After=network.target

[Service]
Type=simple
User=kejarprestasi
Group=kejarprestasi
WorkingDirectory=/var/www/kejarprestasi
ExecStart=/usr/bin/node .node-server/server.node.js
Restart=on-failure
RestartSec=5
StandardOutput=append:/var/www/kejarprestasi/logs/app.log
StandardError=append:/var/www/kejarprestasi/logs/app-err.log

# Environment
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOST=0.0.0.0
EnvironmentFile=/var/www/kejarprestasi/.env

# Hardening (opsional, comment kalau bermasalah)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
